# scripts/drbd-services@.target
[Unit]
Description=Services for DRBD resource %I
BindsTo=drbd-promote@%i.service
