# scripts/drbd-reconfigure-suspend-or-error@.service.in
[Unit]
Description=Reconfigure DRBD %I on-no-(quorum|data) to suspend-io/io-error during start/stop
Documentation=man:drbd-reconfigure-suspend-or-error@.service(7)
PartOf=drbd-services@%i.target

RefuseManualStart=yes
RefuseManualStop=yes

[Service]
SyslogIdentifier=drbd-%I
Type=oneshot
RemainAfterExit=yes
ExecStop=-@sbindir@/drbdsetup resource-options --on-no-quorum io-error --on-no-data io-error %I
ExecStart=-@sbindir@/drbdsetup resource-options --on-no-quorum suspend-io --on-no-data suspend-io %I
